# deploy/techchallenge.service
[Unit]
Description=TechChallengeApp

[Service]
ExecStart=/app/dist/TechChallengeApp serve
Restart=always

[Install]
WantedBy=multi-user.target
