# deploy/techchallenge.service
[Unit]
Description=TechChallengeApp

[Service]
ExecStart=/app/dist/TechChallengeApp serve

[Install]
WantedBy=multi-user.target
